# PoolControlRPi/prerm
#!/bin/bash
systemctl disable poolsmartz
systemctl stop poolsmartz
